# iteration-5/master0/fakeroot/etc/systemd/system/kubelet.service.d/80-nodeip.conf
[Service]
Environment=KUBELET_NODE_IP="172.21.104.21"
